# deploy/robotrader.service
# systemd unit for the RoboTrader engine (PAPER mode).
#   sudo cp deploy/robotrader.service /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now robotrader
#
# LIVE mode is deliberately NOT run under systemd: it requires the typed
# confirmation phrase on an interactive terminal, and stdin here is null —
# a live unit fails safe at startup. Run live inside tmux (see docs/DEPLOY.md).

[Unit]
Description=RoboTrader engine (paper)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=robotrader
WorkingDirectory=/opt/robotrader
# Credentials + alert/heartbeat config; chmod 600, owned by robotrader.
EnvironmentFile=/etc/robotrader.env
ExecStart=/opt/robotrader/.venv/bin/python -m service.engine --config config/paper.yaml
Restart=on-failure
RestartSec=10
# The engine journals every decision; stdout is operational logging only.
StandardOutput=append:/opt/robotrader/logs/engine.out
StandardError=append:/opt/robotrader/logs/engine.err

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/opt/robotrader

[Install]
WantedBy=multi-user.target
